# etc/sv/getty-generic/finish
#!/bin/sh
tty=${PWD##*-}
exit 0

#echo "running utmpset for $tty..." 1>&2
#exec utmpset -w $tty
