# etc/sv/getty-generic/run
#!/bin/sh

tty=${PWD##*-}

[ -r conf ] && . ./conf

GETTY_ARGS="${GETTY_ARGS} -L"

exec chpst -P agetty ${GETTY_ARGS} \
	"${tty}" "${BAUD_RATE}" "${TERM_NAME}"
